# eddiebot/eddiebot_bringup/upstart/eddiebot.conf
description	"bringup eddiebot"

start on net-device-up IFACE=eth1
stop on net-device-down IFACE=eth1

console output
respawn

exec eddiebot-start
pre-stop exec eddiebot-stop
